# kubernetes/kube-proxy.conf
description "Kubernetes Kubelet Service"
start on started kubelet
stop on runlevel [!2345]
respawn

exec hyperkube proxy --v=2 --master="http://master:8080"